# service/phanes-http.service
# phanes-http — systemd unit (design.md Decision 11 EC2 host + Decision 14 deploy).
#
# One-time host install (out of scope for deploy.sh):
#   sudo cp service/phanes-http.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now phanes-http
#
# Adjust User= and the paths to match the EC2 account / repo location.
# deploy.sh then restarts this unit on each ship.

[Unit]
Description=phanes-http — hexa-native discovery SaaS HTTP service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/phanes
# PHANES_HOME lets the service resolve web/static, the job store, jobctl.
Environment=PHANES_HOME=/home/ubuntu/phanes
Environment=PHANES_BIND_HOST=127.0.0.1
Environment=PHANES_BIND_PORT=8787
ExecStart=/home/ubuntu/phanes/bin/phanes-http
Restart=always
RestartSec=2
# graceful: the select-loop accept loop exits on its stop-flag; SIGTERM
# is the default KillSignal and the process closes its listener cleanly.
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
